# roles/common/files/etc/cron.daily/oinkmaster
#!/bin/sh

TMPFILE=$(mktemp)
oinkmaster -C /etc/oinkmaster.conf -o /etc/snort/rules -Q | tee $TMPFILE
test -s $TMPFILE && service suricata restart
rm -f $TMPFILE
