# roles/common/files/etc/cron.d/rkhunter-boot
@reboot	root	/usr/bin/rkhunter --cronjob --report-warnings-only --enable properties
